# deploy/conar255.service
# Template — the setup script replaces __USER__, __HOME__, __APP_DIR__
# before installing to /etc/systemd/system/conar255.service

[Unit]
Description=CONAR Model 255 Advanced Roasting System
After=multi-user.target
Conflicts=getty@tty1.service

[Service]
Type=simple
User=__USER__
Group=__USER__

Environment=HOME=__HOME__

# Start X server with no cursor, running our app directly
ExecStart=/usr/bin/xinit __APP_DIR__/deploy/start-x.local.sh -- :0 vt1 -nocursor -quiet
WorkingDirectory=__APP_DIR__

# Restart on crash (3 second cooldown)
Restart=on-failure
RestartSec=3

# Give access to display, input, and serial devices
SupplementaryGroups=video input dialout render tty

# Allow X to access the tty
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
