# deploy/conar255.service
# Template — the setup script replaces __USER__, __HOME__, __APP_DIR__
# before installing to /etc/systemd/system/conar255.service

[Unit]
Description=CONAR Model 255 Advanced Roasting System
After=multi-user.target
Conflicts=getty@tty1.service

[Service]
Type=simple

# xinit runs as root to access the virtual console;
# start-x.local.sh drops to the app user via sudo
User=root
Environment=HOME=__HOME__
Environment=SDL_AUDIODRIVER=dummy

ExecStart=/usr/bin/xinit __APP_DIR__/deploy/start-x.local.sh -- :0 vt1 -nocursor -quiet
WorkingDirectory=__APP_DIR__

# Restart on crash (3 second cooldown)
Restart=on-failure
RestartSec=3

# Keep stdout/stderr in journalctl
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
